# jobs/crontab
* * * * * python /jobs/get_nba_game_schedules_job.py
